# registry-garbage-collect.timer.jinja
[Unit]
Description=Gitlab Registry automatic garbage collect

[Timer]
OnCalendar={{ on_calendar }}
Unit=gitlab-registry-garbage-collect.service

[Install]
WantedBy=timers.target
